# home/config/systemd/user/restart@.service
[Unit]
Description=Restart systemd services from .path or .timer

[Service]
Type=oneshot
ExecStart=systemctl --user restart %i.service
