# home/config/systemd/user/run@.service
[Unit]
Description=Run anything in ~/bin

[Service]
Type=oneshot
ExecStart=%h/bin/%i
CacheDirectory=%i
WorkingDirectory=%C/%i
