# home/config/systemd/user/waybar.path
[Unit]
Description=Restart waybar on config changes
PartOf=niri.service

[Path]
PathChanged=%h/config/waybar
Unit=restart@waybar.service

[Install]
WantedBy=niri.service
