# home/config/systemd/user/zcompdump.timer
[Unit]
Description=daily Zsh completion dump update

[Timer]
OnCalendar=daily
Unit=run@zcompdump.service

[Install]
WantedBy=timers.target
